# dist/shalture.cron.example
#!/bin/sh
#
# crontab script for shalture
#

# Change this to the directory where you installed it
prefix=$HOME/shalture/
pidfile=var/shalture.pid

cd "$prefix" || {
	echo "$0: couldn't cd to $prefix"
	exit 1
}

if [ -f "$pidfile" ] ; then
  PID=`cat "$pidfile"`

  if kill -0 $PID >/dev/null 2>&1; then
    exit 0
  fi

  rm -f "$pidfile"
fi

echo "$0: couldn't find shalture running, restarting..."

./bin/shalture-services

exit 0
